# packages/cloud-builder/cloud_builder-1.0.6.tar.gz/cloud_builder-1.0.6/systemd/cb-fetch-once.service
[Unit]
Description=CB Project Fetch Once
After=syslog.target
Requires=network-online.target

[Service]
EnvironmentFile=-/etc/cloud_builder
Type=oneshot
ExecStart=/bin/bash -c 'sudo cb-fetch --project "${CB_PROJECT}" --branch "${CB_BRANCH}" --single-shot'

[Install]
WantedBy=multi-user.target
